# config/soip.service.in
[Unit]
Description=SDR OFDM IP Tunnel - Wireless point-to-point link
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=@SOIP_BIN_PATH@ --config @SOIP_CONFIG_PATH@
WorkingDirectory=@INSTALL_DIR@
User=@SERVICE_USER@
Group=@SERVICE_USER@
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW
DeviceAllow=/dev/net/tun rw
SupplementaryGroups=netdev

[Install]
WantedBy=multi-user.target
